# source/systemd/dmi-spark-dns@.timer
[Unit]
Description=dmi-spark DNS refresh

[Timer]
OnActiveSec=%i
OnUnitActiveSec=%i
AccuracySec=1s
Unit=dmi-spark-dns.service
